# deploy/systemd/recipehud-backend.service
[Unit]
Description=Recipe HUD backend
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/opt/recipehud
# Let wlopm reach the labwc compositor socket (Pi OS autologin runs as uid 1000).
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=RECIPEHUD_DB_PATH=/opt/recipehud/data/recipehud.db
ExecStart=/opt/recipehud/.venv/bin/uvicorn recipehud.main:app --host 0.0.0.0 --port 8000 --app-dir /opt/recipehud/backend
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
